# init.d/firstrun
#!/sbin/runscript
#
# Copyright 2005 TUBITAK/UEKAE
# Licensed under the GNU General Public License, version 2.
# See the file http://www.gnu.org/copyleft/gpl.txt.

start() {
        ebegin "First run service starting, please be patient..."

        # 05gcc fix
        echo "PATH=\"/usr/i686-pc-linux-gnu/gcc-bin/3.4.4\"" > /etc/env.d/05gcc
        echo "ROOTPATH=\"/usr/i686-pc-linux-gnu/gcc-bin/3.4.4\"" >> /etc/env.d/05gcc
        echo "MANPATH=\"/usr/share/gcc-data/i686-pc-linux-gnu/3.4.4/man\"" >> /etc/env.d/05gcc
        echo "INFOPATH=\"/usr/share/gcc-data/i686-pc-linux-gnu/3.4.4/info\"" >> /etc/env.d/05gcc
        echo "LDPATH=\"/usr/lib/gcc/i686-pc-linux-gnu/3.4.4\"" >> /etc/env.d/05gcc

        # ld.so.conf fix
        echo "/usr/lib" > /etc/ld.so.conf
        echo "/usr/lib/opengl/xorg-x11/lib" >> /etc/ld.so.conf
        echo "/usr/qt/3/lib"  >> /etc/ld.so.conf
        echo "/usr/kde/3.4/lib"  >> /etc/ld.so.conf
        echo "/usr/lib/gcc/i686-pc-linux-gnu/3.4.4/" >> /etc/ld.so.conf

        ldconfig
        /usr/sbin/update-environ.py

        # fix for some zip problems (see bug #622)
        chmod 1777 /tmp
        chmod 1777 /var/tmp
        chmod 700 /root/
        
        # Run postinstall scripts
        pisi configure-pending

        # configure xorg
        /sbin/xorg.py
        # Remove itself
        rc-update del firstrun

        # add needed servises
        rc-update add alsasound boot
        rc-update add dbus default
        rc-update add hald default 
        rc-update add xdm default

    # Calculate deptree
        depscan.sh
        eend $?
}
