# init.d/firstrun
#!/sbin/runscript
#
# Copyright 2005 TUBITAK/UEKAE
# Licensed under the GNU General Public License, version 2.
# See the file http://www.gnu.org/copyleft/gpl.txt.

depend() {
    after comar
}

start() {
        ebegin "First run service starting, please be patient..."

        /usr/bin/update-environ.py
        ldconfig

        # Run postinstall scripts
        export LC_ALL=tr_TR.UTF-8
        export LANG=tr_TR.UTF-8
        pisi configure-pending 2>&1 &> /root/pisi.log

        # configure xorg
        /sbin/xorg.py
        # Remove itself
        rc-update del firstrun

        # add needed servises
        rc-update add alsasound boot
        rc-update add dbus default
        rc-update add hald default 
        rc-update add xdm default

        # Calculate deptree
        depscan.sh
        eend $?
}
