# init.d/net.lo
#!/sbin/runscript
# Copyright © 2005  TUBITAK/UEKAE
# Licensed under the GNU General Public License, version 2.
# See the file http://www.gnu.org/copyleft/gpl.txt.
#
# Original work belongs Gentoo Linux

depend() {
	need checkroot
}

start() {
    ifconfig lo 127.0.0.1 up 2>/dev/null
    /sbin/route add -net 127.0.0.0 netmask 255.0.0.0 gw 127.0.0.1 dev lo 2> /dev/null
}

# vim:ts=4
